# services/jarvis/jarvis-max-access.service
[Unit]
Description=Nexa Access API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/jarvis-max-access
EnvironmentFile=/opt/jarvis-max-access/jarvis-max-access.env
ExecStart=/opt/jarvis-max-access/venv/bin/uvicorn jarvis_max_access:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
